# supervisors/drivers/pad_drivers.conf
; Port Agent supervisord config file
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9010        ; (ip_address:port specifier, *:port for all iface)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
childlogdir=%(here)s/supervisor_logs

[supervisorctl]
serverurl=unix://supervisor.sock ; use a unix:// URL  for a unix socket

[program:RS01SLBS-LJ01A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SLBS-LJ01A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SUM1-LJ01B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SUM1-LJ01B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:CE04OSBP-LJ01C]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver CE04OSBP-LJ01C %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:CE02SHBP-LJ01D]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver CE02SHBP-LJ01D %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXBS-LJ03A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXBS-LJ03A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SBPS-LV01A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SBPS-LV01A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SUM1-LV01B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SUM1-LV01B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:CE04OSBP-LV01C]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver CE04OSBP-LV01C %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXPS-LV03A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXPS-LV03A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SLBS-MJ01A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SLBS-MJ01A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SUM1-MJ01B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SUM1-MJ01B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:CE02SHBP-MJ01C]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver CE02SHBP-MJ01C %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXBS-MJ03A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXBS-MJ03A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03ASHS-MJ03B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03ASHS-MJ03B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03INT1-MJ03C]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03INT1-MJ03C %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03INT2-MJ03D]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03INT2-MJ03D %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03ECAL-MJ03E]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03ECAL-MJ03E %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03CCAL-MJ03F]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03CCAL-MJ03F %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SLBS-PN01A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SLBS-PN01A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SHDR-PN01B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SHDR-PN01B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01OSBP-PN01C]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01OSBP-PN01C %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SHBP-PN01D]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SHBP-PN01D %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXBS-PN03A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXBS-PN03A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXSM-PN03B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXSM-PN03B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-OTNNA]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-OTNNA %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-OTNSA]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-OTNSA %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-F10NA]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-F10NA %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-F10SA]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-F10SA %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-PFE00]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-PFE00 %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-TS00N]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-TS00N %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-TS00S]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-TS00S %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-UPS0A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-UPS0A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:SSRSPACC-UPS0B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver SSRSPACC-UPS0B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SBPS-SC01A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SBPS-SC01A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01OSBP-SC01B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01OSBP-SC01B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXPS-SC03A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXPS-SC03A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SBPS-PC01A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SBPS-PC01A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01OSBP-PC01B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01OSBP-PC01B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXPS-PC03A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXPS-PC03A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SBPS-SF01A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SBPS-SF01A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01OSBP-SF01B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01OSBP-SF01B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXPS-SF03A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXPS-SF03A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SLBS-PD01A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SLBS-PD01A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS01SHDR-PD01B]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS01SHDR-PD01B %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:RS03AXBS-PD03A]
directory=%(here)s/driver_logs/%(program_name)s
command=run_driver mi.platform.rsn.rsn_platform_driver RSNPlatformDriver RS03AXBS-PD03A %(ENV_EVENTS_QUEUE)s %(ENV_PARTICLES_QUEUE)s
autostart=false

[program:oms_extractor]
directory=%(here)s/driver_logs/%(program_name)s
command=oms_extractor oms_extract_config.yml
autostart=true
autorestart=true
