# supervisors/pa1/pa1.conf
; Port Agent supervisord config file
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9009        ; (ip_address:port specifier, *:port for all iface)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
environment=AMQP_URL="amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production"

[supervisorctl]
serverurl=unix://supervisor.sock ; use a unix:// URL  for a unix socket

[program:pa_stats_shovel]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" port_agent_stats port_agent_stats qpid://uft20.ooi.rutgers.edu port_agent_stats
autostart=true
autorestart=true

[program:events_shovel]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" ha.instrument_events ha.instrument_events qpid://ooiufs03.ooi.rutgers.edu Ingest.instrument_events
autostart=true
autorestart=true

[program:oms_alerts_shovel]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" ha.oms_alerts ha.instrument_events qpid://ooiufs03.ooi.rutgers.edu Ingest.instrument_oms_events
autostart=true
autorestart=true

[program:particles_shovel]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" ha.instrument_particles ha.instrument_particles qpid://ooiufs03.ooi.rutgers.edu Ingest.instrument_particles
autostart=true
autorestart=true

[program:particles_shovel_uft4]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" uft4_instrument_particles ha.instrument_particles qpid://uframe-4-test.ooi.rutgers.edu Ingest.instrument_particles
autostart=false
autorestart=false

[program:pa_monitor]
command=python -m ooi_port_agent.tools.monitor monitor.conf
autostart=true
autorestart=true
redirect_stderr=true

[program:RS01SLBS-LJ01A-12-CTDPFB101]
command=port_agent rsn 10.33.5.12 2101 2102 RS01SLBS-LJ01A-12-CTDPFB101
autostart=true
autorestart=true

[program:RS01SLBS-LJ01A-12-DOSTAD101]
command=port_agent rsn 10.33.5.12 2101 2102 RS01SLBS-LJ01A-12-DOSTAD101
autostart=true
autorestart=true

[program:CE04OSBP-LJ01C-05-ADCPSI103]
command=port_agent rsn 10.33.12.5 2101 2102 CE04OSBP-LJ01C-05-ADCPSI103
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-4F-PCO2WA101]
command=port_agent rsn 10.33.3.209 2101 2102 RS01SBPS-SF01A-4F-PCO2WA101
autostart=true
autorestart=true

[program:RS01SBPS-PC01A-05-ADCPTD102]
command=port_agent rsn 10.33.3.143 2101 2102 RS01SBPS-PC01A-05-ADCPTD102
autostart=true
autorestart=true

[program:RS01SBPS-PC01A-4A-CTDPFA103]
command=port_agent rsn 10.33.3.135 2101 2102 RS01SBPS-PC01A-4A-CTDPFA103
autostart=true
autorestart=true

[program:RS01SBPS-PC01A-4A-DOSTAD103]
command=port_agent rsn 10.33.3.135 2101 2102 RS01SBPS-PC01A-4A-DOSTAD103
autostart=true
autorestart=true

[program:CE02SHBP-LJ01D-08-OPTAAD106]
command=port_agent rsn 10.33.14.8 2101 2102 CE02SHBP-LJ01D-08-OPTAAD106
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-4A-NUTNRA101]
command=port_agent rsn 10.33.3.206 2101 2102 RS01SBPS-SF01A-4A-NUTNRA101
autostart=true
autorestart=true

[program:CE02SHBP-LJ01D-09-PCO2WB103]
command=port_agent rsn 10.33.14.9 2101 2102 CE02SHBP-LJ01D-09-PCO2WB103
autostart=true
autorestart=true

[program:RS01SBPS-PC01A-07-CAMDSC102]
command=port_agent camds 10.33.3.146 8889 '/san_data/' RS01SBPS-PC01A-07-CAMDSC102
autostart=false
autorestart=true

[program:RS01SBPS-PC01A-4B-PHSENA102]
command=port_agent rsn 10.33.3.136 2101 2102 RS01SBPS-PC01A-4B-PHSENA102
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-4B-VELPTD102]
command=port_agent rsn 10.33.3.207 2101 2102 RS01SBPS-SF01A-4B-VELPTD102
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-4B-VELPTD106]
command=port_agent rsn 10.33.10.207 2101 2102 CE04OSPS-SF01B-4B-VELPTD106
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-3B-OPTAAD101]
command=port_agent rsn 10.33.3.202 2101 2102 RS01SBPS-SF01A-3B-OPTAAD101
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-3B-OPTAAD105]
command=port_agent rsn 10.33.10.202 2101 2102 CE04OSPS-SF01B-3B-OPTAAD105
autostart=true
autorestart=true

[program:CE02SHBP-LJ01D-07-VEL3DC108]
command=port_agent rsn 10.33.14.7 2101 2102 CE02SHBP-LJ01D-07-VEL3DC108
autostart=true
autorestart=true

[program:RS01SLBS-LJ01A-11-OPTAAC103]
command=port_agent rsn 10.33.5.11 2101 2102 RS01SLBS-LJ01A-11-OPTAAC103
autostart=true
autorestart=true

[program:CE04OSBP-LJ01C-07-VEL3DC107]
command=port_agent rsn 10.33.12.7 2101 2102 CE04OSBP-LJ01C-07-VEL3DC107
autostart=true
autorestart=true

[program:CE04OSBP-LJ01C-06-CTDBPO108]
command=port_agent rsn 10.33.12.6 2101 2102 CE04OSBP-LJ01C-06-CTDBPO108
autostart=true
autorestart=true

[program:CE04OSBP-LJ01C-06-DOSTAD108]
command=port_agent rsn 10.33.12.6 2101 2102 CE04OSBP-LJ01C-06-DOSTAD108
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-3D-SPKIRA101]
command=port_agent rsn 10.33.3.204 2101 2102 RS01SBPS-SF01A-3D-SPKIRA101
autostart=true
autorestart=true

[program:CE04OSPS-PC01B-4A-CTDPFA109]
command=port_agent rsn 10.33.10.135 2101 2102 CE04OSPS-PC01B-4A-CTDPFA109
autostart=true
autorestart=true

[program:CE04OSPS-PC01B-4A-DOSTAD109]
command=port_agent rsn 10.33.10.135 2101 2102 CE04OSPS-PC01B-4A-DOSTAD109
autostart=true
autorestart=true

[program:RS01SLBS-LJ01A-10-ADCPTE101]
command=port_agent rsn 10.33.5.10 2101 2102 RS01SLBS-LJ01A-10-ADCPTE101
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-3C-PARADA101]
command=port_agent rsn 10.33.3.203 2101 2102 RS01SBPS-SF01A-3C-PARADA101
autostart=true
autorestart=true

[program:CE04OSPS-PC01B-05-ZPLSCB102]
command=port_agent zplsc 10.33.10.143 21 'data/Platform_OPS' '/san_data/' ooi 994ef22 CE04OSPS-PC01B-05-ZPLSCB102
autostart=true
autorestart=true

[program:RS01SLBS-LJ01A-05-HPIESA101]
command=port_agent rsn 10.33.5.5 2101 2102 RS01SLBS-LJ01A-05-HPIESA101
autostart=true
autorestart=true

[program:RS01SLBS-MJ01A-12-VEL3DB101]
command=port_agent rsn 10.33.2.12 2101 2102 RS01SLBS-MJ01A-12-VEL3DB101
autostart=true
autorestart=true

[program:CE04OSPS-PC01B-4B-PHSENA106]
command=port_agent rsn 10.33.10.136 2101 2102 CE04OSPS-PC01B-4B-PHSENA106
autostart=true
autorestart=true

[program:CE04OSBP-LJ01C-08-OPTAAC104]
command=port_agent rsn 10.33.12.8 2101 2102 CE04OSBP-LJ01C-08-OPTAAC104
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-3D-SPKIRA102]
command=port_agent rsn 10.33.10.204 2101 2102 CE04OSPS-SF01B-3D-SPKIRA102
autostart=true
autorestart=true

[program:CE04OSBP-LJ01C-09-PCO2WB104]
command=port_agent rsn 10.33.12.9 2101 2102 CE04OSBP-LJ01C-09-PCO2WB104
autostart=true
autorestart=true

[program:CE02SHBP-MJ01C-08-CAMDSB107]
command=port_agent camds 10.33.13.8 8889 '/san_data/' CE02SHBP-MJ01C-08-CAMDSB107
autostart=true
autorestart=true

[program:RS01SLBS-MJ01A-06-PRESTA101]
command=port_agent rsn 10.33.2.6 2101 2102 RS01SLBS-MJ01A-06-PRESTA101
autostart=true
autorestart=true

[program:RS01SBPS-PC01A-4C-FLORDD103]
command=port_agent rsn 10.33.3.137 2101 2102 RS01SBPS-PC01A-4C-FLORDD103
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-3C-PARADA102]
command=port_agent rsn 10.33.10.203 2101 2102 CE04OSPS-SF01B-3C-PARADA102
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-2A-CTDPFA107]
command=port_agent rsn 10.33.10.195 2101 2102 CE04OSPS-SF01B-2A-CTDPFA107
autostart=true
autorestart=true

[program:CE02SHBP-MJ01C-07-ZPLSCB101]
command=port_agent zplsc 10.33.13.7 21 'data/Shelf_OPS' '/san_data/' ooi 994ef22 CE02SHBP-MJ01C-07-ZPLSCB101
autostart=true
autorestart=true

[program:RS01SBPS-PC01A-06-VADCPA101_5]
command=port_agent rsn 10.33.3.145 2101 2102 RS01SBPS-PC01A-06-VADCPA101-5
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-3A-FLORTD101]
command=port_agent rsn 10.33.3.201 2101 2102 RS01SBPS-SF01A-3A-FLORTD101
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-2D-PHSENA101]
command=port_agent rsn 10.33.3.197 2101 2102 RS01SBPS-SF01A-2D-PHSENA101
autostart=true
autorestart=true

[program:CE04OSPS-PC01B-4C-PCO2WA105]
command=port_agent rsn 10.33.10.138 2101 2102 CE04OSPS-PC01B-4C-PCO2WA105
autostart=true
autorestart=true

[program:RS01SBPS-PC01A-06-VADCPA101]
command=port_agent rsn 10.33.3.144 2101 2102 RS01SBPS-PC01A-06-VADCPA101-4
autostart=true
autorestart=true

[program:CE02SHBP-LJ01D-10-PHSEND103]
command=port_agent rsn 10.33.14.10 2101 2102 CE02SHBP-LJ01D-10-PHSEND103
autostart=true
autorestart=true

[program:CE02SHBP-LJ01D-06-CTDBPN106]
command=port_agent rsn 10.33.14.6 2101 2102 CE02SHBP-LJ01D-06-CTDBPN106
autostart=true
autorestart=true

[program:CE02SHBP-LJ01D-06-DOSTAD106]
command=port_agent rsn 10.33.14.6 2101 2102 CE02SHBP-LJ01D-06-DOSTAD106
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-3A-FLORTD104]
command=port_agent rsn 10.33.10.201 2101 2102 CE04OSPS-SF01B-3A-FLORTD104
autostart=true
autorestart=true

[program:CE04OSBP-LJ01C-10-PHSEND107]
command=port_agent rsn 10.33.12.10 2101 2102 CE04OSBP-LJ01C-10-PHSEND107
autostart=true
autorestart=true

[program:CE02SHBP-LJ01D-05-ADCPTB104]
command=port_agent rsn 10.33.14.5 2101 2102 CE02SHBP-LJ01D-05-ADCPTB104
autostart=true
autorestart=true

[program:CE04OSBP-LV01C-06-CAMDSB106]
command=port_agent camds 10.33.9.6 8889 '/san_data/' CE04OSBP-LV01C-06-CAMDSB106
autostart=true
autorestart=true

[program:RS01SBPS-SF01A-2A-CTDPFA102]
command=port_agent rsn 10.33.3.195 2101 2102 RS01SBPS-SF01A-2A-CTDPFA102
autostart=true
autorestart=true
