# supervisors/antelope/antelope.conf
; Port Agent supervisord config file
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9009        ; (ip_address:port specifier, *:port for all iface)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]

[supervisorctl]
serverurl=unix://supervisor.sock ; use a unix:// URL  for a unix socket

[program:PA-RS01SLBS-LJ01A-09-HYDBBA102]
command=port_agent antelope localhost 7000 RS01SLBS-LJ01A-09-HYDBBA102
autostart=true

[program:PA-RS01SBPS-PC01A-08-HYDBBA103]
command=port_agent antelope localhost 7000 RS01SBPS-PC01A-08-HYDBBA103
autostart=true

[program:PA-CE04OSBP-LJ01C-11-HYDBBA105]
command=port_agent antelope localhost 7000 CE04OSBP-LJ01C-11-HYDBBA105
autostart=true

[program:PA-CE02SHBP-LJ01D-11-HYDBBA106]
command=port_agent antelope localhost 7000 CE02SHBP-LJ01D-11-HYDBBA106
autostart=true

[program:PA-RS03AXBS-LJ03A-09-HYDBBA302]
command=port_agent antelope localhost 7000 RS03AXBS-LJ03A-09-HYDBBA302
autostart=true

[program:PA-RS03AXPS-PC03A-08-HYDBBA303]
command=port_agent antelope localhost 7000 RS03AXPS-PC03A-08-HYDBBA303
autostart=true

[program:PA-RS01SLBS-MJ01A-05-HYDLFA101]
command=port_agent antelope localhost 7100 RS01SLBS-MJ01A-05-HYDLFA101
autostart=true

[program:PA-RS01SUM1-LJ01B-05-HYDLFA104]
command=port_agent antelope localhost 7100 RS01SUM1-LJ01B-05-HYDLFA104
autostart=true

[program:PA-RS03AXBS-MJ03A-05-HYDLFA301]
command=port_agent antelope localhost 7100 RS03AXBS-MJ03A-05-HYDLFA301
autostart=true

[program:PA-RS03ECAL-MJ03E-09-HYDLFA304]
command=port_agent antelope localhost 7100 RS03ECAL-MJ03E-09-HYDLFA304
autostart=true

[program:PA-RS03CCAL-MJ03F-06-HYDLFA305]
command=port_agent antelope localhost 7100 RS03CCAL-MJ03F-06-HYDLFA305
autostart=true

[program:PA-RS01SUM1-LJ01B-05-OBSBBA101]
command=port_agent antelope localhost 7100 RS01SUM1-LJ01B-05-OBSBBA101
autostart=true

[program:PA-RS01SLBS-MJ01A-05-OBSBBA102]
command=port_agent antelope localhost 7100 RS01SLBS-MJ01A-05-OBSBBA102
autostart=true

[program:PA-RS03CCAL-MJ03F-06-OBSBBA301]
command=port_agent antelope localhost 7100 RS03CCAL-MJ03F-06-OBSBBA301
autostart=true

[program:PA-RS03ECAL-MJ03E-09-OBSBBA302]
command=port_agent antelope localhost 7100 RS03ECAL-MJ03E-09-OBSBBA302
autostart=true

[program:PA-RS03AXBS-MJ03A-05-OBSBBA303]
command=port_agent antelope localhost 7100 RS03AXBS-MJ03A-05-OBSBBA303
autostart=true

[program:PA-RS01SUM1-LJ01B-08-OBSSPA101]
command=port_agent antelope localhost 7100 RS01SUM1-LJ01B-08-OBSSPA101
autostart=true

[program:PA-RS01SUM1-LJ01B-07-OBSSPA102]
command=port_agent antelope localhost 7100 RS01SUM1-LJ01B-07-OBSSPA102
autostart=true

[program:PA-RS01SUM1-LJ01B-06-OBSSPA103]
command=port_agent antelope localhost 7100 RS01SUM1-LJ01B-06-OBSSPA103
autostart=true

[program:PA-RS03ASHS-MJ03B-06-OBSSPA301]
command=port_agent antelope localhost 7100 RS03ASHS-MJ03B-06-OBSSPA301
autostart=true

[program:PA-RS03ASHS-MJ03B-05-OBSSPA302]
command=port_agent antelope localhost 7100 RS03ASHS-MJ03B-05-OBSSPA302
autostart=true

[program:PA-RS03ECAL-MJ03E-05-OBSSPA303]
command=port_agent antelope localhost 7100 RS03ECAL-MJ03E-05-OBSSPA303
autostart=true

[program:PA-RS03ECAL-MJ03E-08-OBSSPA304]
command=port_agent antelope localhost 7100 RS03ECAL-MJ03E-08-OBSSPA304
autostart=true

[program:PA-RS03INT2-MJ03D-05-OBSSPA305]
command=port_agent antelope localhost 7100 RS03INT2-MJ03D-05-OBSSPA305
autostart=true
