# supervisors/pa2/pa2.conf
; Port Agent supervisord config file
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9009        ; (ip_address:port specifier, *:port for all iface)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
environment=AMQP_URL="amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production"

[supervisorctl]
serverurl=unix://supervisor.sock ; use a unix:// URL  for a unix socket

[program:uframe_3_test_shovel]
command=shovel amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production uft3_instrument_particles ha.instrument_particles qpid://uframe-3-test.ooi.rutgers.edu Ingest.instrument_particles
autostart=false
autorestart=false

[program:uframe_2_test_shovel]
command=shovel amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production uft2_instrument_particles ha.instrument_particles qpid://uframe-2-test.ooi.rutgers.edu Ingest.instrument_particles
autostart=false
autorestart=false

[program:uframe_test_shovel]
command=shovel amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production uft1_instrument_particles ha.instrument_particles qpid://uframe-test.ooi.rutgers.edu Ingest.instrument_particles
autostart=false
autorestart=false

[program:pa_stats_shovel]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" port_agent_stats port_agent_stats qpid://uft20.ooi.rutgers.edu port_agent_stats
autostart=false
autorestart=false

[program:events_shovel]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" ha.instrument_events ha.instrument_events qpid://ooiufs03.ooi.rutgers.edu Ingest.instrument_events
autostart=false
autorestart=false

[program:particles_shovel]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" ha.instrument_particles ha.instrument_particles qpid://ooiufs03.ooi.rutgers.edu Ingest.instrument_particles
autostart=false
autorestart=false

[program:particles_shovel_uft20]
command=shovel "amqp://ooi-production:T1cZp68WanMZ53dQ@localhost/production" uft20_instrument_particles ha.instrument_particles qpid://uft20.ooi.rutgers.edu Ingest.instrument_particles
autostart=false
autorestart=false

[program:pa_monitor]
command=python -m ooi_port_agent.tools.monitor monitor.conf
autostart=true
redirect_stderr=true

[program:RS03ASHS-MJ03B-07-TMPSFA301]
command=port_agent rsn 10.31.7.7 2101 2102 RS03ASHS-MJ03B-07-TMPSFA301
autostart=true
autorestart=true

[program:RS03ASHS-MJ03B-09-BOTPTA304]
command=port_agent botpt 10.31.7.9 9338 9337 RS03ASHS-MJ03B-09-BOTPTA304
autostart=true
autorestart=true

[program:RS03ASHS-MJ03B-10-CTDPFB304]
command=port_agent rsn 10.31.7.10 2101 2102 RS03ASHS-MJ03B-10-CTDPFB304
autostart=true
autorestart=true

[program:RS03AXBS-LJ03A-05-HPIESA301]
command=port_agent rsn 10.31.5.5 2101 2102 RS03AXBS-LJ03A-05-HPIESA301
autostart=true
autorestart=true

[program:RS03AXPS-PC03A-4B-PHSENA302]
command=port_agent rsn 10.31.3.136 2101 2102 RS03AXPS-PC03A-4B-PHSENA302
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-3C-PARADA301]
command=port_agent rsn 10.31.3.203 2101 2102 RS03AXPS-SF03A-3C-PARADA301
autostart=true
autorestart=true

[program:RS03CCAL-MJ03F-05-BOTPTA301]
command=port_agent botpt 10.31.6.5 9338 9337 RS03CCAL-MJ03F-05-BOTPTA301
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-09-THSPHA301]
command=port_agent rsn 10.31.8.9 2101 2102 RS03INT1-MJ03C-09-THSPHA301
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-06-MASSPA301-RGA]
command=port_agent tcp 10.31.8.6 4001 RS03INT1-MJ03C-06-MASSPA301-RGA
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-10-TRHPHA301]
command=port_agent rsn 10.31.8.10 2101 2102 RS03INT1-MJ03C-10-TRHPHA301
autostart=true
autorestart=true

[program:RS03AXPS-PC03A-06-VADCPA301]
command=port_agent rsn 10.31.3.144 2101 2102 RS03AXPS-PC03A-06-VADCPA301-4
autostart=true
autorestart=true

[program:RS01SUM2-MJ01B-12-ADCPSK101]
command=port_agent rsn 10.33.7.12 2101 2102 RS01SUM2-MJ01B-12-ADCPSK101
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-4A-NUTNRA301]
command=port_agent rsn 10.31.3.206 2101 2102 RS03AXPS-SF03A-4A-NUTNRA301
autostart=true
autorestart=true

[program:RS01SUM2-MJ01B-05-CAMDSB103]
command=port_agent camds 10.33.7.5 8889 '/san_data/' RS01SUM2-MJ01B-05-CAMDSB103
autostart=true
autorestart=true

[program:RS03ECAL-MJ03E-06-BOTPTA302]
command=port_agent botpt 10.31.10.6 9338 9337 RS03ECAL-MJ03E-06-BOTPTA302
autostart=true
autorestart=true

[program:RS03AXBS-LJ03A-12-CTDPFB301]
command=port_agent rsn 10.31.5.12 2101 2102 RS03AXBS-LJ03A-12-CTDPFB301
autostart=true
autorestart=true

[program:RS03AXBS-LJ03A-12-DOSTAD301]
command=port_agent rsn 10.31.5.12 2101 2102 RS03AXBS-LJ03A-12-DOSTAD301
autostart=true
autorestart=true

[program:RS03AXPS-PC03A-07-CAMDSC302]
command=port_agent camds 10.31.3.146 8889 '/san_data/' RS03AXPS-PC03A-07-CAMDSC302
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-07-D1000A301]
command=port_agent tcp 10.31.8.7 4002 RS03INT1-MJ03C-07-D1000A301
autostart=true
autorestart=true

[program:RS03INT2-MJ03D-12-VEL3DB304]
command=port_agent rsn 10.31.9.12 2101 2102 RS03INT2-MJ03D-12-VEL3DB304
autostart=true
autorestart=true

[program:RS03AXPS-PC03A-05-ADCPTD302]
command=port_agent rsn 10.31.3.143 2101 2102 RS03AXPS-PC03A-05-ADCPTD302
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-07-RASFLA301]
command=port_agent tcp 10.31.8.7 4001 RS03INT1-MJ03C-07-RASFLA301
autostart=true
autorestart=true

[program:RS03AXPS-PC03A-4A-CTDPFA303]
command=port_agent rsn 10.31.3.135 2101 2102 RS03AXPS-PC03A-4A-CTDPFA303
autostart=true
autorestart=true

[program:RS03AXPS-PC03A-4A-DOSTAD303]
command=port_agent rsn 10.31.3.135 2101 2102 RS03AXPS-PC03A-4A-DOSTAD303
autostart=true
autorestart=true

[program:RS01SUM2-MJ01B-06-MASSPA101-RGA]
command=port_agent tcp 10.33.7.6 4001 RS01SUM2-MJ01B-06-MASSPA101-RGA
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-2A-CTDPFA302]
command=port_agent rsn 10.31.3.195 2101 2102 RS03AXPS-SF03A-2A-CTDPFA302
autostart=true
autorestart=true

[program:RS03INT2-MJ03D-06-BOTPTA303]
command=port_agent botpt 10.31.9.6 9338 9337 RS03INT2-MJ03D-06-BOTPTA303
autostart=true
autorestart=true

[program:RS01SUM2-MJ01B-06-MASSPA101-TURBO]
command=port_agent tcp 10.33.7.6 4002 RS01SUM2-MJ01B-06-MASSPA101-TURBO
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-3A-FLORTD301]
command=port_agent rsn 10.31.3.201 2101 2102 RS03AXPS-SF03A-3A-FLORTD301
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-05-CAMDSB303]
command=port_agent camds 10.31.8.5 8889 '/san_data/' RS03INT1-MJ03C-05-CAMDSB303
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-06-MASSPA301-MCU]
command=port_agent tcp 10.31.8.6 4003 RS03INT1-MJ03C-06-MASSPA301-MCU
autostart=true
autorestart=true

[program:RS03AXBS-MJ03A-06-PRESTA301]
command=port_agent rsn 10.31.2.6 2101 2102 RS03AXBS-MJ03A-06-PRESTA301
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-4F-PCO2WA301]
command=port_agent rsn 10.31.3.209 2101 2102 RS03AXPS-SF03A-4F-PCO2WA301
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-3D-SPKIRA301]
command=port_agent rsn 10.31.3.204 2101 2102 RS03AXPS-SF03A-3D-SPKIRA301
autostart=true
autorestart=true

[program:RS01SUM2-MJ01B-06-MASSPA101-MCU]
command=port_agent tcp 10.33.7.6 4003 RS01SUM2-MJ01B-06-MASSPA101-MCU
autostart=true
autorestart=true

[program:RS01SUM1-LJ01B-12-VEL3DB104]
command=port_agent rsn 10.33.8.12 2101 2102 RS01SUM1-LJ01B-12-VEL3DB104
autostart=true
autorestart=true

[program:RS03AXBS-LJ03A-10-ADCPTE303]
command=port_agent rsn 10.31.5.10 2101 2102 RS03AXBS-LJ03A-10-ADCPTE303
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-4B-VELPTD302]
command=port_agent rsn 10.31.3.207 2101 2102 RS03AXPS-SF03A-4B-VELPTD302
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-2D-PHSENA301]
command=port_agent rsn 10.31.3.197 2101 2102 RS03AXPS-SF03A-2D-PHSENA301
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-07-PPSDNA301]
command=port_agent tcp 10.31.8.7 4003 RS03INT1-MJ03C-07-PPSDNA301
autostart=true
autorestart=true

[program:RS03INT1-MJ03C-06-MASSPA301-TURBO]
command=port_agent tcp 10.31.8.6 4002 RS03INT1-MJ03C-06-MASSPA301-TURBO
autostart=true
autorestart=true

[program:RS03AXPS-PC03A-06-VADCPA301_5]
command=port_agent rsn 10.31.3.145 2101 2102 RS03AXPS-PC03A-06-VADCPA301-5
autostart=true
autorestart=true

[program:RS01SUM1-LJ01B-09-PRESTB102]
command=port_agent rsn 10.33.8.9 2101 2102 RS01SUM1-LJ01B-09-PRESTB102
autostart=true
autorestart=true

[program:RS03AXPS-SF03A-3B-OPTAAD301]
command=port_agent rsn 10.31.3.202 2101 2102 RS03AXPS-SF03A-3B-OPTAAD301
autostart=true
autorestart=true

[program:RS03AXPS-PC03A-4C-FLORDD303]
command=port_agent rsn 10.31.3.137 2101 2102 RS03AXPS-PC03A-4C-FLORDD303
autostart=true
autorestart=true

[program:RS03AXBS-MJ03A-12-VEL3DB301]
command=port_agent rsn 10.31.2.12 2101 2102 RS03AXBS-MJ03A-12-VEL3DB301
autostart=true
autorestart=true

[program:RS03AXBS-LJ03A-11-OPTAAC303]
command=port_agent rsn 10.31.5.11 2101 2102 RS03AXBS-LJ03A-11-OPTAAC303
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-4A-NUTNRA102]
command=port_agent rsn 10.33.10.206 2101 2102 CE04OSPS-SF01B-4A-NUTNRA102
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-4F-PCO2WA102]
command=port_agent rsn 10.33.10.209 2101 2102 CE04OSPS-SF01B-4F-PCO2WA102
autostart=true
autorestart=true

[program:CE04OSPS-SF01B-2B-PHSENA108]
command=port_agent rsn 10.33.10.197 2101 2102 CE04OSPS-SF01B-2B-PHSENA108
autostart=true
autorestart=true
